# .config/systemd/mihomo.service
[Service]
Restart=always
ExecStart=mihomo
[Install]
WantedBy=default.target
