# config/systemd/mihomo.service
[Service]
Restart=always
ExecStart=%h/.config/mihomo/main
[Install]
WantedBy=default.target
