# config/systemd/mihomo.service
[Unit]
After=network.target
[Service]
Restart=always
ExecStart=%h/.config/mihomo/mihomo
[Install]
WantedBy=default.target
